# workshop/socket-activation/systemd.units/socket-activation-python.service
[Unit]
Description=Testing python socket server

[Service]
WorkingDirectory=/opt/workshop
ExecStart=socket-activation/server